# web20.conf
[program:web21]
command=/usr/local/bin/gunicorn wsgi --bind localhost:4200  
directory=/root/web9-lC21
autostart=true
autorestart=true
